# crontab
0 0 * * * /usr/bin/systemctl restart datalogApp

