# crontab
0 0 * * * /bin/systemctl restart cameraApp
