# services/carApp.service
[Unit]
Description=UI service
After=multi-user.target
[Service]
ExecStart=/usr/bin/python /root/ha/carApp.py
Restart=on-abort
[Install]
WantedBy=multi-user.target
